# qa/qa_api/tools/supervisor-spug.ini
[program:qa-api]
command = bash /data/qa/qa_api/tools/start-api.sh
autostart = true
stdout_logfile = /data/qa/qa_api/logs/api.log
redirect_stderr = true
